# installer/linux/claude-standup.service
[Unit]
Description=Claude Standup Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/claude-standup daemon run
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
