# installer/macos/com.claude-standup.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-standup.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/claude-standup</string>
        <string>daemon</string>
        <string>run</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>~/.claude-standup/daemon.log</string>
    <key>StandardErrorPath</key>
    <string>~/.claude-standup/daemon.log</string>
</dict>
</plist>
